# scripts/cec-panasonic-ampli.service
[Unit]
Description=CEC Panasonic Ampli for Homebridge IR Amplifier
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/homebridge/node_modules/homebridge-ir-amplifier
ExecStart=/var/lib/homebridge/node_modules/homebridge-ir-amplifier/scripts/cec-panasonic-ampli.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
