# scripts/cec-panasonic-ampli.service
[Unit]
Description=CEC Panasonic Ampli for Homebridge IR Amplifier
After=multi-user.target
Wants=homebridge.service

[Service]
Type=simple
ExecStart=/usr/local/bin/cec-panasonic-ampli.sh
Restart=always
RestartSec=15
User=root
StandardOutput=journal
StandardError=journal

# Permissions pour accéder au device CEC
SupplementaryGroups=video

[Install]
WantedBy=multi-user.target
